# cmder_config/profile.d/aliases.sh
# common alias
alias ll='ls -ltra'
alias la='ls -la'
alias pwdw='cygpath -ml `pwd`'
alias docs='cd $HOMEDRIVE/Documents'
alias gf='git fetch --all'
alias gs='git status'