# cmder_config/profile.d/aliases.sh
# common alias
alias ll='ls -ltra'
alias pwdw='cygpath -ml `pwd`'
alias docs='cd $HOMEDRIVE/Documents'