# cmder_config/profile.d/aliases.sh
# common alias
alias ll='ls -ltra'
alias pwdw='cygpath -ml `pwd`'
alias docs='cd $HOMEDRIVE/Documents'

# ikti server
alias ikti.srv='ssh bambang_ebis@10.126.12.212 -i /c/Users/zeruel11/.ssh/iktisrv_id'